# bootLed.service
[Unit]
Description=Boot LED Animation - Amber breathing while system starts
After=basic.target
Before=audioGuestBook.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/admin/rotary-phone-audio-guestbook/src/bootLed.py
Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
